# config/systemd/user/ib-box-spread-nats.service
[Unit]
Description=IB Box Spread - NATS server
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.config/systemd/user/ib-box-spread.env
WorkingDirectory=%h/.config/systemd/user
ExecStart=/bin/bash -c 'cd "${PROJECT_ROOT}" && exec nats-server -js -DV'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
