# config/systemd/user/ib-box-spread-risk_free_rate.service
[Unit]
Description=IB Box Spread - Risk-free rate backend service
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.config/systemd/user/ib-box-spread.env
WorkingDirectory=%h/.config/systemd/user
ExecStart=/bin/bash -c 'cd "${PROJECT_ROOT}/python" && exec python -m uvicorn integration.risk_free_rate_service:app --host 0.0.0.0 --port 8004'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
