# config/systemd/user/ib-box-spread-web.service
[Unit]
Description=IB Box Spread - Web (Vite dev server)
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.config/systemd/user/ib-box-spread.env
WorkingDirectory=%h/.config/systemd/user
ExecStart=/bin/bash -c 'cd "${PROJECT_ROOT}/web" && exec npm run dev -- --host 0.0.0.0 --port 5173'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
