# config/systemd/user/ib-box-spread-tradestation.service
[Unit]
Description=IB Box Spread - TradeStation backend service
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.config/systemd/user/ib-box-spread.env
WorkingDirectory=%h/.config/systemd/user
ExecStart=/bin/bash -c 'cd "${PROJECT_ROOT}/python" && exec python -m uvicorn integration.tradestation_service:app --host 0.0.0.0 --port 8001'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
