# web/scripts/systemd/pwa-rust-backend.service
[Unit]
Description=PWA Rust Backend Service
After=network.target pwa-nats.service
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/ib_box_spread_full_universal
ExecStart=%h/ib_box_spread_full_universal/scripts/start_rust_backend.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pwa-rust-backend
Environment="HOME=%h"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
