# web/scripts/systemd/pwa-jupyterlab.service
[Unit]
Description=PWA JupyterLab Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/ib_box_spread_full_universal
ExecStart=%h/ib_box_spread_full_universal/scripts/run-jupyterlab-service.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pwa-jupyterlab
Environment="HOME=%h"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
