# web/scripts/systemd/pwa-ib-gateway.service
[Unit]
Description=PWA IB Gateway Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/ib_box_spread_full_universal/ib-gateway
# Try run-gateway-with-reload.sh first, then fallback to run-gateway.sh, then bin/run.sh
ExecStart=/bin/bash -c 'cd "%h/ib_box_spread_full_universal/ib-gateway" && if [ -f "run-gateway-with-reload.sh" ]; then exec bash run-gateway-with-reload.sh; elif [ -f "run-gateway.sh" ]; then exec bash run-gateway.sh; elif [ -f "bin/run.sh" ]; then CONFIG_FILE="root/conf.yaml"; [ ! -f "$CONFIG_FILE" ] && CONFIG_FILE="root/conf.tws.yaml"; exec bash bin/run.sh "$CONFIG_FILE"; else echo "Error: No gateway script found"; exit 1; fi'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pwa-ib-gateway
Environment="HOME=%h"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
