# web/scripts/systemd/pwa-discount-bank.service
[Unit]
Description=PWA Discount Bank Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/ib_box_spread_full_universal
ExecStart=%h/ib_box_spread_full_universal/web/scripts/run-discount-bank-service.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pwa-discount-bank
Environment="HOME=%h"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
