# web/scripts/systemd/pwa-ib.service
[Unit]
Description=PWA IB (Interactive Brokers) Service
After=network.target pwa-ib-gateway.service
Wants=network-online.target
Requires=pwa-ib-gateway.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/ib_box_spread_full_universal
ExecStart=%h/ib_box_spread_full_universal/web/scripts/run-ib-service.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pwa-ib
Environment="HOME=%h"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
