# web/scripts/systemd/pwa-nats.service
[Unit]
Description=PWA NATS Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/ib_box_spread_full_universal
ExecStart=/usr/bin/nats-server -c %h/ib_box_spread_full_universal/config/nats-server.conf
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pwa-nats
Environment="HOME=%h"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
