# config/supervisord.conf
; Process manager config for all backend services (Python supervisord).
; Run from project root with PROJECT_ROOT set:
;   export PROJECT_ROOT=/path/to/ib_box_spread_full_universal
;   supervisord -c config/supervisord.conf
;
; Or use the wrapper: ./scripts/run_supervisord.sh
;
; Install supervisord: pip install supervisord  or  apt install supervisor

[unix_http_server]
file = /tmp/supervisor-backend-services.sock

[supervisord]
logfile = /tmp/supervisord-backend-services.log
logfile_maxbytes = 10MB
logfile_backups = 2
loglevel = info
pidfile = /tmp/supervisord-backend-services.pid
nodaemon = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor-backend-services.sock

; ---- Backend HTTP services (must set PROJECT_ROOT before starting) ----

[program:ib]
command = python -m uvicorn integration.ib_service:app --host 0.0.0.0 --port 8002
directory = %(ENV_PROJECT_ROOT)s/python
environment = PYTHONPATH="%(ENV_PROJECT_ROOT)s/python"
autostart = true
autorestart = true
stdout_logfile = /tmp/supervisor-ib.log
stderr_logfile = /tmp/supervisor-ib.err.log

[program:health_dashboard]
command = python -m uvicorn services.health_dashboard:app --host 0.0.0.0 --port 8011
directory = %(ENV_PROJECT_ROOT)s/python
environment = PYTHONPATH="%(ENV_PROJECT_ROOT)s/python"
autostart = true
autorestart = true
stdout_logfile = /tmp/supervisor-health_dashboard.log
stderr_logfile = /tmp/supervisor-health_dashboard.err.log

[program:nats]
command = nats-server -js -DV
directory = %(ENV_PROJECT_ROOT)s
autostart = true
autorestart = true
stdout_logfile = /tmp/supervisor-nats.log
stderr_logfile = /tmp/supervisor-nats.err.log

; Rust backend (optional; comment out if not built)
[program:rust_backend]
command = cargo run --release -- --rest-port 8080 --grpc-port 50051
directory = %(ENV_PROJECT_ROOT)s/agents/backend
autostart = false
autorestart = true
stdout_logfile = /tmp/supervisor-rust_backend.log
stderr_logfile = /tmp/supervisor-rust_backend.err.log
